# ansible/roles/app/templates/django_app_supervisor.conf.j2
[program:django_app]
command={{ django_base_dir }}/django_app_gunicorn_start.sh
autorestart=true
directory={{ django_base_dir }}
autostart=true
autorestart=true